# backups/dotfiles/udev/90-keyboard-backlight.rules
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="tpacpi::kbd_backlight", RUN+="/bin/chmod 666 /sys%p/brightness"